# supervisor-dev.conf
; Reference config. In production this lives on the server (paths under
; /var/www/FastForward-Express) and IS used by supervisorctl; in dev workers are
; normally started by hand. Keep this in sync with the prod config.
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
directory=/home/brandon/Desktop/FastForward-Express
; Launch via the wrapper, which blocks until MySQL is reachable before starting
; queue:work. This keeps each process alive past startsecs during a transient DB
; outage so the group can never be driven to FATAL and silently stay dead.
command=/home/brandon/Desktop/FastForward-Express/tools/run-worker.sh
numprocs=2
user=brandon
autostart=true
autorestart=true
startsecs=5
; High retry ceiling is a backstop only; the launcher's DB wait is the real fix.
startretries=999
stopsignal=QUIT
killasgroup=true
redirect_stderr=true
stdout_logfile=/home/brandon/Desktop/FastForward-Express/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; Event listener for health monitoring
events=PROCESS_STATE,PROCESS_STATE_STARTING,PROCESS_STATE_RUNNING,PROCESS_STATE_BACKOFF,PROCESS_STATE_STOPPING,PROCESS_STATE_EXITED,PROCESS_STATE_STOPPED,PROCESS_STATE_FATAL

[eventlistener:worker-state-monitor]
command=/home/brandon/Desktop/FastForward-Express/resources/assets/bash/proc_events.sh
events=PROCESS_STATE
buffer_size=100
environment=WEBHOOK_URL="http://127.0.0.1:8000/api/webhooks/supervisor/status",WEBHOOK_TOKEN="dev-token"
